# etc/udev/rules.d/99-g903-reset.rules
# Nutzt den robusten Modalias-String fuer das G903 Kabel (C091)
ACTION=="remove", SUBSYSTEM=="usb", ENV{MODALIAS}=="usb:v046DpC091*", RUN+="/usr/bin/systemd-run --service-type=oneshot /usr/local/bin/g903-reset.sh"
