# etc/systemd/system/rclone-mount@.service
[Unit]
Description=Rclone Mount für %i
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/rclone mount %i: /mnt/%i \
  --config /etc/rclone/rclone.conf \
  --allow-other \
  --vfs-cache-mode full \
  --vfs-cache-max-size 20G \
  --vfs-cache-max-age 24h \
  --vfs-read-chunk-size 64M \
  --vfs-read-chunk-size-limit 2G \
  --dir-cache-time 1000h \
  --poll-interval 30s
ExecStop=/usr/bin/fusermount -u /mnt/%i
Restart=on-failure
RestartSec=10
User=kevin

[Install]
WantedBy=multi-user.target


